# app/public/week39.json
{"version":"2026.06.21.1","module":{"id":"week39","title":"Threat Intelligence & MITRE ATT&CK","durationText":"Week 39 (Days 1-2)","focus":"Map adversary behavior to MITRE ATT&CK, consume threat intelligence, and operationalize IOCs","output":"Ability to map detections to ATT&CK techniques, consume threat feeds, and understand the intelligence lifecycle","topics":[{"id":"we39d01","title":"MITRE ATT&CK Framework","description":"Navigate the ATT&CK matrix, map detected activity to techniques, and use the Navigator","type":"learn","duration":"4 hours","content":":::objectives\n- Navigate the 14 ATT&CK tactics and their relationships\n- Map real attack behaviors to specific techniques and sub-techniques\n- Use the ATT&CK Navigator to visualize detection coverage\n- Identify common techniques used in real-world attacks\n:::\n\n:::concept\nMITRE ATT&CK is a knowledge base of adversary tactics and techniques based on real-world observations.\n:::\n\n### The 14 ATT&CK Tactics\n\n| # | Tactic | Example Technique |\n|---|--------|-------------------|\n| 1 | **Reconnaissance** | T1595 - Active Scanning |\n| 2 | **Resource Development** | T1583 - Acquire Infrastructure |\n| 3 | **Initial Access** | T1566 - Phishing |\n| 4 | **Execution** | T1059 - Command Interpreter |\n| 5 | **Persistence** | T1053 - Scheduled Task |\n| 6 | **Privilege Escalation** | T1068 - Exploitation for Priv Esc |\n| 7 | **Defense Evasion** | T1027 - Obfuscated Files |\n| 8 | **Credential Access** | T1003 - OS Credential Dumping |\n| 9 | **Discovery** | T1087 - Account Discovery |\n| 10 | **Lateral Movement** | T1021 - Remote Services |\n| 11 | **Collection** | T1005 - Data from Local System |\n| 12 | **Command and Control** | T1071 - Application Layer Protocol |\n| 13 | **Exfiltration** | T1041 - Exfil Over C2 |\n| 14 | **Impact** | T1486 - Data Encrypted (Ransomware) |\n\n### Key Techniques\n\n**T1566 - Phishing:**\n```\nDetection: Email gateway (attachments .iso, .img, .vhd), endpoint (OUTLOOK.EXE -> cmd.exe)\n```\n\n**T1059 - Command and Scripting Interpreter:**\n```\nT1059.001 - PowerShell: -enc, IEX, Invoke-Expression\nT1059.003 - cmd.exe: certutil, bitsadmin\nT1059.005 - VBScript: cscript.exe, wscript.exe\nDetection: Sysmon Event 1, Windows Event 4688, PowerShell ScriptBlock Logging\n```\n\n**T1053 - Scheduled Task (Persistence):**\n```\nDetection: Windows Event 4698, Sysmon Event 1 (schtasks /create)\n```\n\n**T1021 - Remote Services (Lateral Movement):**\n```\nT1021.001 - RDP: Logon Type 10\nT1021.002 - SMB: net use, psexec\nT1021.004 - SSH\nDetection: Event 4624 Type 10, Event 5145, Sysmon Event 3\n```\n\n### ATT&CK Navigator\n\n1. Go to https://mitre-attack.github.io/attack-navigator/\n2. Create a layer and color cells:\n   - **Green** = Detection exists\n   - **Yellow** = Partial detection\n   - **Red** = No detection (gap)\n   - **Gray** = Not applicable\n\n### Common Attack Chains\n\n**Ransomware (Conti-style):**\n```\nT1566 -> T1059.001 -> T1003 -> T1021.002 -> T1053 -> T1486 -> T1567\n```\n\n:::classwork\n**Lab:** Open ATT&CK Navigator, create a layer for your detection capabilities, color by coverage, identify top 3 gaps, write a Sigma rule for one gap.\n:::","aiPrompt":"","labUrl":"","labTitle":"","interviewQuestion":"How do you map a specific attack behavior to an ATT&CK technique?","interviewAnswer":"Identify the adversary goal (e.g., stealing credentials), find the corresponding tactic (Credential Access) and technique (T1003). This mapping ensures systematic detection coverage, reveals gaps, and enables comparison against threat actor profiles.","quiz":[{"question":"Which ATT&CK tactic includes \"Phishing\" (T1566)?","options":["Initial Access","Execution","Persistence","Reconnaissance"],"correctAnswerIndex":0,"difficulty":"beginner","explanation":"Phishing is an Initial Access technique.","certTags":["MITRE ATT&CK"]},{"question":"What is the difference between a tactic and a technique?","options":["Tactics are goals (why), techniques are methods (how)","Same thing","Tactics for blue team, techniques for red","Techniques are broader"],"correctAnswerIndex":0,"difficulty":"beginner","explanation":"Tactics represent the goal; techniques represent the method.","certTags":["MITRE ATT&CK"]},{"question":"Which technique covers PowerShell encoded commands?","options":["T1059.001 - PowerShell","T1059.003 - cmd.exe","T1204 - User Execution","T1205 - Service Execution"],"correctAnswerIndex":0,"difficulty":"intermediate","explanation":"T1059.001 specifically covers PowerShell scripting interpreter.","certTags":["MITRE ATT&CK"]},{"question":"What does the ATT&CK Navigator visualize?","options":["Detection coverage across the ATT&CK matrix","Geographic attacker locations","Network topology","SIEM log rates"],"correctAnswerIndex":0,"difficulty":"beginner","explanation":"The Navigator layers detection capabilities over the ATT&CK matrix.","certTags":["MITRE ATT&CK"]},{"question":"T1021.001 corresponds to which lateral movement method?","options":["RDP","SMB","SSH","WMI"],"correctAnswerIndex":0,"difficulty":"intermediate","explanation":"T1021.001 is Remote Desktop Protocol.","certTags":["MITRE ATT&CK"]},{"question":"Which tactic covers \"Data Encrypted for Impact\" (T1486)?","options":["Impact","Exfiltration","Collection","Defense Evasion"],"correctAnswerIndex":0,"difficulty":"beginner","explanation":"T1486 is ransomware encryption under the Impact tactic.","certTags":["MITRE ATT&CK"]},{"question":"What does \"Defense Evasion\" (TA0005) cover?","options":["Techniques to avoid detection - obfuscation, disabling defenses, injection","Credential stealing","Lateral movement","Persistence"],"correctAnswerIndex":0,"difficulty":"intermediate","explanation":"Defense Evasion includes T1027 (Obfuscation), T1562 (Impair Defenses), T1055 (Injection).","certTags":["MITRE ATT&CK"]},{"question":"In a ransomware chain, which technique comes FIRST?","options":["T1566 (Phishing) - initial access","T1486 (Encryption) - impact","T1003 (Credential Dump)","T1021 (Lateral Movement)"],"correctAnswerIndex":0,"difficulty":"intermediate","explanation":"The kill chain begins with Initial Access (Phishing).","certTags":["MITRE ATT&CK","Ransomware"]},{"question":"Which technique does Windows Event 4698 detect?","options":["T1053 (Scheduled Task) - task creation","T1059 (Command Interpreter)","T1003 (Credential Dump)","T1486 (Data Encrypted)"],"correctAnswerIndex":0,"difficulty":"intermediate","explanation":"Event 4698 records Scheduled Task creation.","certTags":["MITRE ATT&CK","Windows"]},{"question":"Which sub-technique of T1003 dumps LSASS memory?","options":["T1003.001 - LSASS Memory","T1003.002 - SAM","T1003.003 - NTDS","T1003.006 - DCSync"],"correctAnswerIndex":0,"difficulty":"advanced","explanation":"T1003.001 specifically covers LSASS memory dumping.","certTags":["MITRE ATT&CK","Credential Access"]},{"question":"What color represents partial detection in the Navigator?","options":["Yellow","Green","Red","Blue"],"correctAnswerIndex":0,"difficulty":"beginner","explanation":"Green = full coverage, Yellow = partial, Red = gap.","certTags":["MITRE ATT&CK"]},{"question":"Which tactic includes T1078 (Valid Accounts)?","options":["Initial Access and Persistence","Credential Access only","Defense Evasion only","Lateral Movement only"],"correctAnswerIndex":0,"difficulty":"advanced","explanation":"T1078 spans Initial Access, Persistence, Privilege Escalation, and Defense Evasion.","certTags":["MITRE ATT&CK"]},{"question":"Why is ATT&CK important for SOC teams?","options":["Common language for attacks, gap analysis, coverage measurement","Replaces SIEM tools","Only for red teams","Compliance framework like PCI-DSS"],"correctAnswerIndex":0,"difficulty":"beginner","explanation":"ATT&CK provides standardized vocabulary and helps identify detection gaps.","certTags":["MITRE ATT&CK","SOC"]},{"question":"What is the difference between T1566.001 and T1566.002?","options":["001 = malicious attachment; 002 = URL in email body","001 = mass phishing; 002 = targeted","001 = Windows; 002 = Linux","No difference"],"correctAnswerIndex":0,"difficulty":"intermediate","explanation":"T1566.001 delivers via file attachment; T1566.002 via URL.","certTags":["MITRE ATT&CK","Phishing"]},{"question":"Which technique covers creating a new user account for persistence?","options":["T1136 - Create Account","T1098 - Account Manipulation","T1078 - Valid Accounts","T1531 - Account Access Removal"],"correctAnswerIndex":0,"difficulty":"intermediate","explanation":"T1136 covers new account creation for persistence.","certTags":["MITRE ATT&CK"]}]},{"id":"we39d02","title":"Threat Intelligence & IOC Operationalization","description":"Understand TLP, consume threat feeds, use STIX/TAXII, and operationalize IOCs in detection","type":"practice","duration":"4 hours","content":":::objectives\n- Apply Traffic Light Protocol (TLP)\n- Use IOC types and formats (STIX/TAXII) in detection\n- Consume threat feeds from OTX, abuse.ch, and VirusTotal\n- Operationalize IOCs into Wazuh/Sigma detection rules\n:::\n\n### Traffic Light Protocol (TLP)\n\n| TLP | Name | Sharing Scope |\n|-----|------|---------------|\n| **TLP:RED** | Restricted | Only within your organization |\n| **TLP:AMBER** | Amber | Your org + clients with need-to-know |\n| **TLP:GREEN** | Green | Within your community (ISACs) |\n| **TLP:CLEAR** | Clear | Unrestricted - public sharing |\n\n:::warning\nNever share TLP:RED information outside your organization.\n:::\n\n### IOC Types\n\n| Type | Example | Use In Detection |\n|------|---------|------------------|\n| IPv4 | 10.0.0.50 | Firewall blocks, SIEM watchlists |\n| Domain | evil-payload.com | DNS blocklists |\n| URL | http://evil.com/payload.exe | Proxy rules |\n| File Hash (SHA256) | e3b0c44298fc... | EDR scanning |\n| JA3 Hash | e7d705a3286e... | TLS fingerprint detection |\n\n### STIX/TAXII\n\n**STIX** is a format for encoding threat intel:\n\n```json\n{\n  \"type\": \"indicator\",\n  \"id\": \"indicator--a1b2c3d4-e5f6-7890-abcd-ef1234567890\",\n  \"name\": \"Cobalt Strike C2 Server\",\n  \"pattern\": \"[ipv4-addr:value = '10.0.0.50']\",\n  \"confidence\": 85\n}\n```\n\n**TAXII** is the protocol for sharing STIX data via HTTP/REST APIs.\n\n### Threat Feed Sources\n\n```bash\n# OTX AlienVault\ncurl -X GET \"https://otx.alienvault.com/api/v1/indicators/IPv4/10.0.0.50/pulses\"   -H \"X-OTX-API-KEY: YOUR_KEY\"\n\n# ThreatFox (abuse.ch)\ncurl -X POST \"https://threatfox-api.abuse.ch/api/v1/\"   -d '{\"query\": \"search_ioc\", \"search_term\": \"10.0.0.50\"}'\n\n# MalwareBazaar (abuse.ch)\ncurl -X POST \"https://mb-api.abuse.ch/api/v1/\"   -d '{\"query\": \"get_info\", \"hash\": \"d41d8cd98f00b204e9800998ecf8427e\"}'\n```\n\n### Operationalizing IOCs in Wazuh\n\n```bash\n# Download threat feed\ncurl -o /var/ossec/etc/lists/abusech-ip.txt   \"https://urlhaus-api.abuse.ch/v1/urls/recent/100/\"\n```\n\n```xml\n<list>etc/lists/abusech-ip.txt</list>\n\n<rule id=\"100300\" level=\"10\">\n  <if_sid>5712</if_sid>\n  <list field=\"srcip\">etc/lists/abusech-ip.txt</list>\n  <description>Connection from known malicious IP (abuse.ch)</description>\n  <group>threat_intel,pci_dss_11.4</group>\n</rule>\n```\n\n### Intelligence Lifecycle\n\n```\nDirection -> Collection -> Processing -> Analysis -> Dissemination -> Feedback\n```\n\n:::classwork\n**Lab:** Create OTX account, find a ransomware pulse, extract 5 IOCs, verify on VirusTotal, create Wazuh CDB list and rule.\n:::","aiPrompt":"","labUrl":"","labTitle":"","interviewQuestion":"Explain TLP:AMBER vs TLP:GREEN.","interviewAnswer":"TLP:AMBER restricts sharing to your org and specific clients. TLP:GREEN allows sharing within your community (ISACs, partners). AMBER is bilateral, GREEN is community-wide.","quiz":[{"question":"What does TLP:RED mean?","options":["Share only within your org - need-to-know","Share with clients and partners","Share within industry group","Share publicly"],"correctAnswerIndex":0,"difficulty":"beginner","explanation":"TLP:RED is the most restrictive sharing level.","certTags":["Threat Intel"]},{"question":"What is STIX?","options":["Standardized format for encoding threat intel data","Protocol for real-time sharing","Firewall rule type","SIEM aggregation format"],"correctAnswerIndex":0,"difficulty":"beginner","explanation":"STIX is a language for describing threat intel in machine-readable format.","certTags":["Threat Intel"]},{"question":"What protocol does TAXII use?","options":["HTTP/REST API","SMTP","SNMP","FTP"],"correctAnswerIndex":0,"difficulty":"intermediate","explanation":"TAXII uses RESTful HTTP APIs to share STIX objects.","certTags":["Threat Intel"]},{"question":"Which feed provides malware hashes and samples?","options":["MalwareBazaar (abuse.ch)","OTX AlienVault","AbuseIPDB","WHOIS"],"correctAnswerIndex":0,"difficulty":"intermediate","explanation":"MalwareBazaar provides SHA256 hashes, file types, and tags.","certTags":["Threat Intel"]},{"question":"Correct order of the Intelligence Lifecycle?","options":["Direction -> Collection -> Processing -> Analysis -> Dissemination -> Feedback","Collection -> Direction -> Analysis -> Processing","Processing -> Collection -> Direction","Analysis -> Collection -> Processing -> Direction"],"correctAnswerIndex":0,"difficulty":"beginner","explanation":"The lifecycle starts with Direction (what questions need answering).","certTags":["Threat Intel"]},{"question":"How do you operationalize IOCs in Wazuh?","options":["Create CDB list and reference in detection rule via <list> tag","Paste IPs into <white_list>","Import CSV via dashboard","Wazuh cannot use external lists"],"correctAnswerIndex":0,"difficulty":"intermediate","explanation":"Wazuh uses CDB lists referenced in rules via the <list> element.","certTags":["Wazuh","Threat Intel"]},{"question":"What does a JA3 hash represent?","options":["TLS client fingerprint - identifies application making connections","File hash","DNS query signature","Email header hash"],"correctAnswerIndex":0,"difficulty":"advanced","explanation":"JA3 hashes TLS Client Hello fields to identify the client application.","certTags":["Threat Intel","Network"]},{"question":"Which TLP level allows sharing within industry ISACs?","options":["TLP:GREEN","TLP:AMBER","TLP:RED","TLP:CLEAR"],"correctAnswerIndex":0,"difficulty":"beginner","explanation":"TLP:GREEN permits community-wide sharing.","certTags":["Threat Intel"]},{"question":"What is the difference between a threat feed and a pulse?","options":["Feed = continuous IOC stream; pulse = single report grouping related IOCs","Same thing","Pulse = real-time API; feed = static file","Feed = IPs only; pulse = hashes only"],"correctAnswerIndex":0,"difficulty":"intermediate","explanation":"A feed is ongoing; a pulse groups IOCs for a specific campaign.","certTags":["Threat Intel"]},{"question":"What does \"Direction\" in the Intelligence Lifecycle involve?","options":["Define what questions need answering based on threat landscape","Direct traffic to threat servers","Classify intel by direction","Point SIEM to threat feed"],"correctAnswerIndex":0,"difficulty":"intermediate","explanation":"Direction is planning: what threats does the org face?","certTags":["Threat Intel"]},{"question":"STIX confidence: 85 represents what level?","options":["High confidence - likely accurate from multiple sources","Low - unverified","Medium - needs validation","Absolute certainty"],"correctAnswerIndex":0,"difficulty":"advanced","explanation":"85 is high confidence, based on multiple sources or direct observation.","certTags":["Threat Intel","STIX"]},{"question":"Which abuse.ch service provides active C2 server IPs?","options":["ThreatFox","MalwareBazaar","URLhaus","Feodo Tracker"],"correctAnswerIndex":0,"difficulty":"intermediate","explanation":"ThreatFox provides IOCs for active malware C2 infrastructure.","certTags":["Threat Intel"]},{"question":"What is the \"diamond model\" in threat intel?","options":["Framework linking adversary, capability, infrastructure, victim","Encryption method","SIEM dashboard layout","IOC confidence scoring"],"correctAnswerIndex":0,"difficulty":"advanced","explanation":"The Diamond Model maps four vertices: Adversary, Capability, Infrastructure, Victim.","certTags":["Threat Intel"]},{"question":"Why age out old IOCs from detection lists?","options":["Old IOCs become stale - infrastructure changes, IPs reassigned, cause FPs","Saves disk space","Old IOCs are always FP","SIEM performance degrades"],"correctAnswerIndex":0,"difficulty":"intermediate","explanation":"Threat infrastructure changes frequently; stale IOCs match benign traffic.","certTags":["Threat Intel"]},{"question":"How to verify a feed IOC is still active before blocking?","options":["Cross-reference VirusTotal, check WHOIS expiry, test connectivity","Add immediately if in feed","Wait 30 days","Only use 7+ day old IOCs"],"correctAnswerIndex":0,"difficulty":"advanced","explanation":"Verify IOCs through multiple sources to reduce false positives.","certTags":["Threat Intel"]}]}]}}